# deploy/systemd/commerceops-certbot-renew.timer
[Unit]
Description=Run CommerceOps certificate renewal twice daily

[Timer]
OnCalendar=*-*-* 03,15:15:00
RandomizedDelaySec=30m
Persistent=true
Unit=commerceops-certbot-renew.service

[Install]
WantedBy=timers.target
